# infra/roles/deploy_tunnel/templates/cloudflared.service.j2
[Unit]
Description=Cloudflare Tunnel
After=network.target

[Service]
TimeoutStartSec=0
Type=notify
User=cloudflared
ExecStart=/usr/local/bin/cloudflared tunnel run --token {{ cloudflare_tunnel_token }}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
